# dedupe.timer
# https://wiki.gbe0.com/en/linux/snippets/btrfs-deduplicate

[Unit]
Description=Deduplicate BTRFS directories weekly
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=7200
Unit=dedupe.service

[Install]
WantedBy=timers.target
